# deploy/systemd/btx-start-network-agent.service
[Unit]
Description=BTX Start network observer agent
After=network-online.target btxd.service
Wants=network-online.target

[Service]
Type=simple
User=btxstart
Group=btxstart
WorkingDirectory=/opt/btx-start
EnvironmentFile=/etc/btx-start/backend.env
ExecStart=/opt/btx-start/.venv/bin/btxstart-network-agent
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/btx-start

[Install]
WantedBy=multi-user.target
